# clonedeploy_pbi/ports/www/pyload/files/pyload.in
#!/bin/sh
#
# PROVIDE: %%PORTNAME%%
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# %%PORTNAME%%_enable (bool):	Set to NO by default.
#			Set it to YES to enable it.
# %%PORTNAME%%_user:	The user account %%PORTNAME%% daemon runs as what
#			you want it to be. It uses '%%USER%%' user by
#			default. Do not sets it as empty or it will run
#			as root.
# %%PORTNAME%%_group:	The group account %%PORTNAME%% daemon runs as what
#			you want it to be. It uses '%%GROUP%%' group by
#			default. Do not sets it as empty or it will run
#			as wheel.
# %%PORTNAME%%_data_dir:	Directory where %%PORTNAME%% configuration
#			data is stored.
#			Default: /var/db/%%PORTNAME%%

. /etc/rc.subr
name="%%PORTNAME%%"
rcvar="${name}_enable"
load_rc_config ${name}

: ${%%PORTNAME%%_enable:="NO"}
: ${%%PORTNAME%%_user:="%%USER%%"}
: ${%%PORTNAME%%_group:="%%GROUP%%"}
: ${%%PORTNAME%%_data_dir:="/var/db/%%PORTNAME%%"}

pidfile="/var/run/%%PORTNAME%%/%%PORTNAME%%.pid"
command="%%PYTHON_CMD%%"
%%PORTNAME%%_flags="%%DATADIR%%/pyload/pyLoadCore.py --configdir=${pyload_data_dir} -p ${pidfile}"
command_args="--daemon"
extra_commands="setup useredit"

start_precmd="%%PORTNAME%%_prestart"
%%PORTNAME%%_prestart() {
	if [ ! -d ${pidfile%/*} ]; then
		install -d -o ${%%PORTNAME%%_user} -g ${%%PORTNAME%%_group} ${pidfile%/*}
	fi

	if [ ! -d ${%%PORTNAME%%_data_dir} ]; then
		install -d -o ${%%PORTNAME%%_user} -g ${%%PORTNAME%%_group} ${%%PORTNAME%%_data_dir}
	fi
}

setup_cmd="%%PORTNAME%%_setup"
%%PORTNAME%%_setup() {
	export PATH=${PATH}:%%PREFIX%%/bin
	${command} ${%%PORTNAME%%_flags} -s
	chown -R ${%%PORTNAME%%_user}:${%%PORTNAME%%_group} ${%%PORTNAME%%_data_dir}
}

useredit_cmd="%%PORTNAME%%_useredit"
%%PORTNAME%%_useredit() {
	${command} ${%%PORTNAME%%_flags} -u
}

run_rc_command "$1"
